# help_tools/rootfs/etc/inittab
# /etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::sysinit:/bin/ash
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

